# app-admin/chkreboot/files/chk_reboot.service
[Unit]
Description=service to check autossh reboot lock

[Service]
ExecStart=/usr/bin/chk_reboot.sh
